# scripts/oraia-notion-sync.timer
[Unit]
Description=ORAIA Notion Sync — every 30 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
