# scripts/oraia-notion-sync.timer
[Unit]
Description=ORAIA Notion Sync — every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
